# dbReports/debian/postrm.in
#!/bin/bash
# Copyright (C) 2010-2012 Ion Torrent Systems, Inc. All Rights Reserved

set -e

case "$1" in
    remove)
    
    #old,legacy log directory
    if [ -d /var/log/anaserve ]; then
    	rm -rf /var/log/anaserve
    fi
    
    # remove sys-V init script links
    DAEMONS=(
        ionCrawler
        ionJobServer
        ionArchive
        ionPlugin
        celeryd
        )
    for DAEMON in ${DAEMONS[@]}; do
        update-rc.d ${DAEMON} remove || true
    done
    ;;
esac

exit 0
