# dbReports/debian/preinst.in
#!/bin/bash
# Copyright (C) 2010 Ion Torrent Systems, Inc. All Rights Reserved
set -e
#set -x
echo "preinst $@"
case "$1" in
	install|upgrade)
    if [ -e	"/etc/init.d/ionCrawler" ]; then
		invoke-rc.d ionCrawler stop
	fi
	if [ -e	"/etc/init.d/ionJobServer" ]; then
		invoke-rc.d ionJobServer stop
	fi
	if [ -e	"/etc/init.d/ionPlugin" ]; then
		invoke-rc.d ionPlugin stop
	fi
	
        invoke-rc.d apache2 stop
    
    # ionArchive is removed as of this version - make sure sysv links are gone
    update-rc.d -f ionArchive remove
	;;
esac

exit 0
