# dbReports/debian/prerm.in
#!/bin/bash
# Copyright (C) 2010 Ion Torrent Systems, Inc. All Rights Reserved
set -e
echo "prerm $@"

SETTINGS=@ION_INSTALL_PREFIX@/iondb/settings.py

case "$1" in
    upgrade|remove)
	if [ -e	"/etc/init.d/ionCrawler" ]; then
		invoke-rc.d ionCrawler stop
	fi
	if [ -e	"/etc/init.d/ionJobServer" ]; then
		invoke-rc.d ionJobServer stop
	fi
	if [ -e	"/etc/init.d/ionPlugin" ]; then
		invoke-rc.d ionPlugin stop
	fi
	if [ -e	"/etc/init.d/celerybeat" ]; then
		invoke-rc.d celeryd stop
	fi
	if [ -e	"/etc/init.d/celeryd" ]; then
		invoke-rc.d celeryd stop
	fi
	if [ -e "/etc/init/DjangoFTP" ]; then
	    invoke-rc.d DjangoFTP stop
	fi

    if [ ! -e @ION_INSTALL_PREFIX@/.computenode ]; then
	    #--- Remove symbolic links	---#
	    rm -f /var/www/output
        rm -rf /var/www/media
        rm -rf /var/www/site_media
        rm -f /results/analysis/output/tmap-*
        rm -f /results/analysis/output/disabled/tmap-*
    fi

	if [ "$1" = "remove" ]; then
		if [ -d @ION_INSTALL_PREFIX@/iondb ]; then
			find "@ION_INSTALL_PREFIX@/iondb" -name "*.log" -exec rm -f {} \;
			rm -f "@ION_INSTALL_PREFIX@/iondb/bin/install_script.py"
		fi
	fi
	if [ -d @ION_INSTALL_PREFIX@/iondb ]; then
		find "@ION_INSTALL_PREFIX@/iondb" -name "*.pyc" -exec rm -f {} \;
	fi

	# Remove torrent-server vhost
	/usr/sbin/a2dissite torrent-server.conf
	rm -f /etc/apache2/sites-enabled/torrent-server

	/usr/sbin/a2ensite 000-default #
	invoke-rc.d apache2 reload || true

    #----------------------------------
    # Remove files created during package install
    #----------------------------------
    rm -f "@ION_INSTALL_PREFIX@/iondb/templates/rundb/php_base.html"
    rm -f "@ION_INSTALL_PREFIX@/iondb/templates/rundb/php_base_plugin.html"

    ;;
esac

exit 0
