# tsconfig/debian/prerm.in
#!/bin/bash
# Copyright (C) 2010 Ion Torrent Systems, Inc. All Rights Reserved

set -e

case "$1" in
    upgrade|remove)
        if [ -e /etc/rc.local ]; then
            sed -i "/.*TSchkconfig.*/d" /etc/rc.local
            sed -i '/qdel.*/d' /etc/rc.local
            sed -i '/clear/d' /etc/rc.local
        fi

        find @PYTHON_SITE_PACKAGES@/ion_tsconfig -name \*.pyc -exec rm {} \; || true
    ;;
esac

exit 0
