# RSM/debian/postinst.in
#!/bin/sh

set -e

# change owner & group since its owned by root due to the install running as root
chown -R ionadmin.www-data /opt/ion/RSM
chmod -R g+w /opt/ion/RSM

# Install the RSM service if we are the master node
if [ -e "/opt/ion/.masternode" ]; then
	# start the service
	if [ -e "/etc/init.d/RSM_Launch" ]; then
		/etc/init.d/RSM_Launch start
	fi

	# see if this was already registered, if not, register it now
	CMD="update-rc.d RSM_Launch defaults 98"
	FOUND=`find /etc/rc?.d -type l | grep RSM_Launch | wc -l`
	if [ "$FOUND" -ne 0 ]; then
		echo "RSM_Launch previously installed"
	else
		$CMD
	fi
fi

